# deploy/supervisor/enaf-scheduler.conf
; Laravel scheduler as a persistent worker (replaces the per-minute `* * * * * schedule:run` cron).
; `schedule:work` stays alive and runs the schedule every minute, which lets sub-minute tasks
; (everyFiveSeconds / everyThirtySeconds) actually fire — the per-minute cron only ticked once a minute.
;
; IMPORTANT: run EITHER this OR the cron, never both, or every scheduled task runs twice.
;
; Install to: /etc/supervisor/conf.d/enaf-scheduler.conf
; Adjust `command` php path (run `which php`) and `user` to match your app before applying.

[program:enaf-scheduler]
process_name=%(program_name)s
command=/usr/bin/php /var/www/enaf/artisan schedule:work
directory=/var/www/enaf
numprocs=1
autostart=true
autorestart=true
startsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/www/enaf/storage/logs/scheduler.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=120
stopsignal=TERM
